# ansible/files/crontab/root
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# backup at 1am every day
1 * * * /srv/signalboost/bin/db/backup
