# ansible/files/crontab/root
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# backup at 5am every day
0 5 * * /srv/signalboost/bin/backup
